# provision/roles/supervisor/files/config.ini
[program:django]
command=/home/{{ app_user }}/current/.virtualenv/bin/gunicorn_django -w 3 -t 60 --settings=settings --pythonpath=.
directory=/home/{{ app_user }}/current/backend
numprocs={{ app_worker_num }}
autostart=true
autorestart=false
stopsignal=INT
stopwaitsecs=2
startsecs=2
redirect_stderr=true
stdout_logfile={{ app_logs_folder }}/gunicorn.log
stdout_logfile_backups=20
stdout_logfile_maxbytes=20MB
user={{ app_user }}
environment=LANG="en_US.UTF-8",USER="{{ app_user }}",HOME="/home/{{ app_user }}"

[program:celery]
command=/home/{{ app_user }}/current/.virtualenv/bin/python manage.py celery worker -l INFO -E
directory=/home/{{ app_user }}/current/backend
numprocs={{ celery_worker_num }}
autostart=true
autorestart=false
stopsignal=INT
stopwaitsecs=2
startsecs=2
redirect_stderr=true
stdout_logfile={{ app_logs_folder }}/celery.log
stdout_logfile_backups=20
stdout_logfile_maxbytes=20MB
user={{ app_user }}
environment=LANG="en_US.UTF-8",USER="{{ app_user }}",HOME="/home/{{ app_user }}"
